# rootfs/etc/s6-overlay/s6-rc.d/init-obsidian-login/up
/etc/s6-overlay/scripts/init-obsidian-login
