# rootfs/etc/s6-overlay/s6-rc.d/init-check-auth/up
/etc/s6-overlay/scripts/init-check-auth
